# fs/etc/services.d/php/run
#!/command/with-contenv sh

exec /usr/sbin/php-fpm8.1 --nodaemonize
