# fs/etc/services.d/nginx/run
#!/command/with-contenv sh
# shellcheck shell=sh

exec nginx -g "daemon off;"
